# rc.local
#!/bin/sh -e

sudo sysctl -p
while ! ps -e | grep honeygain; do true; done && sudo renice -n -20 -p $(pgrep honeygain)

exit 0
