# etc/rc.local
#!/bin/sh -e

sudo modprobe nf_conntrack
sudo modprobe tcp_bbr
sudo modprobe zram
sudo mount --bind /run/systemd/resolve/stub-resolv.conf /etc/resolv.conf
sudo sysctl -p
sudo systemctl restart systemd-networkd.service

exit 0
