# etc/rc.local
#!/bin/sh -e

sudo modprobe nf_conntrack
sudo mount --bind /run/systemd/resolve/stub-resolv.conf /etc/resolv.conf
sudo sysctl -p
sudo systemctl restart networking

exit 0
