# etc/rc.local
#!/bin/sh -e

sudo modprobe nf_conntrack
sudo sysctl -p

exit 0
