# etc/rc.local
#!/bin/sh -e

sudo mount --bind /run/systemd/resolve/stub-resolv.conf /etc/resolv.conf
sudo modprobe nf_conntrack
sudo sysctl -p

exit 0
